# Config Files/beets/crontab
# Crontab for beets automated imports
# This file should be placed in the beets config directory
# and will be picked up by the container

# Run music import every minute (for testing)
* * * * * abc /scripts/import-music.sh >> /config/cron-import.log 2>&1

# Optional: Run daily at 3:00 AM instead (uncomment to use)
# 0 3 * * * abc /scripts/import-music.sh >> /config/cron-import.log 2>&1
